# ansible-pull.service
[Unit]
Description=Ansible Pull for KVM Host Optimization
After=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/kvm-host-optimization
ExecStartPre=/bin/sleep $((RANDOM % 300))
ExecStart=/usr/bin/ansible-pull -U https://github.com/your-repo/kvm-host-optimization.git -i inventory/hosts.yml site.yml
TimeoutSec=1800
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target