# roles/network_nat_persistence/templates/network-setup.service.j2
[Unit]
Description=Network Setup for KVM Host
After=network.target libvirtd.service
Wants=libvirtd.service

[Service]
Type=oneshot
ExecStart=/etc/network/host-optimization/network-setup.sh
RemainAfterExit=yes
TimeoutSec=60

[Install]
WantedBy=multi-user.target