# roles/focused_kvm_optimization/templates/host-cpu-affinity.service.j2
[Unit]
Description=Set CPU affinity for host processes
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/set-host-cpu-affinity.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target